# etc/systemd/system/xconsole.target
[Unit]
Description=XConsole

[Install]
WantedBy=multi-user.target
